# portal_server/supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:api-server]
command=python3 python_webr_fastcgi.py 
user = www-data
autostart=true
